# .config/systemd/user/generate-xresources.service
[Service]
ExecStart=/usr/bin/sh -c 'generate depends xresources | entr -np generate xresources'
[Install]
WantedBy=user-graphical-login.target
