# .config/systemd/user/generate-kitty.service
[Service]
ExecStart=/usr/bin/sh -c 'generate depends kitty | entr -np generate kitty'
[Install]
WantedBy=user-graphical-login.target
