# .config/systemd/user/megasync.service
[Unit]
Description=Megasync

[Service]
ExecStart=/usr/bin/megasync

[Install]
WantedBy=user-graphical-login.target
