# .config/systemd/user/redshift.service
[Unit]
Description=Redshift

[Service]
ExecStart=/usr/bin/redshift

[Install]
WantedBy=user-graphical-login.target
