# .config/systemd/user/redshift.service
[Service]
ExecStart=/usr/bin/redshift
[Install]
WantedBy=user-graphical-login.target
