# .config/systemd/user/sync-rss.timer
[Unit]
Description=Hourly rss sync

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

