# .config/systemd/user/sync-anime.timer
[Unit]
Description=Hourly anime sync

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

