# .config/systemd/user/sync-rss.service
[Unit]
Conflicts/usr/bin/shutdown.target
[Service]
ExecStart=%h/.local/script/sfeed_download_images
