# .config/systemd/user/sync-mail.timer
[Unit]
Description=Hourly mail sync

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

