# .config/systemd/user/sync-anime.service
[Unit]
Description=Sync anime database
Conflicts=shutdown.target

[Service]
ExecStart=/usr/local/bin/anilist fetch
