# .config/systemd/user/volume-bar.service
[Service]
ExecStart=/usr/bin/sh -c 'volume-bar'
[Install]
WantedBy=user-graphical-login.target
