# .config/systemd/user/wallpaper.timer
[Unit]
Description=Minutely switch wallpaper

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

