# .config/systemd/user/xbanish.service
[Unit]
Description=Xbanish

[Service]
ExecStart=/usr/bin/xbanish

[Install]
WantedBy=user-graphical-login.target

