# config/systemd/user/activate-linux.service
[Service]
Type=exec
ExecStart=/usr/bin/activate-linux
Restart=on-failure
RestartSec=5s


[Install]
WantedBy=user-graphical-login.target

