# config/systemd/user/activate-linux.service
[Service]
Type=exec
ExecStart=/usr/bin/activate-linux
[Install]
WantedBy=user-graphical-login.target

