# config/systemd/user/anki.service
[Service]
Type=exec
ExecStart=/sbin/anki

[Unit]
PartOf=graphical-session.target

[Install]
WantedBy=graphical-session.target
