# config/systemd/user/bar-watch.service
[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart bar.service

[Install]
WantedBy=user-graphical-login.target
