# config/systemd/user/bar-watch.path
[Path]
PathChanged=/usr/local/bin/lemonbar-maker
PathChanged=%h/.local/script/bar
PathChanged=%h/.config/xorg/resources
[Install]
WantedBy=user-graphical-login.target
