# config/systemd/user/bar.service
[Service]
Type=exec
ExecStart=/usr/bin/waybar

[Install]
WantedBy=user-graphical-login.target
