# config/systemd/user/bar.service
[Service]
ExecStart=/home/hejsil/.local/script/bar
[Install]
WantedBy=user-graphical-login.target
