# config/systemd/user/bspwm-restart.path
[Path]
PathChanged=%h/.config/bspwm/bspwmrc
PathChanged=%h/.local/script/colors.sh
[Install]
WantedBy=user-graphical-login.target
