# config/systemd/user/caprine.service
[Service]
ExecStart=/usr/bin/caprine
[Install]
WantedBy=user-graphical-login.target

