# config/systemd/user/discord.service
[Service]
Type=exec
ExecStart=/usr/bin/sh -c 'sleep 10s && /usr/bin/vesktop'
[Install]
WantedBy=user-graphical-login.target
