# config/systemd/user/discord.service
[Service]
ExecStart=/usr/bin/discord
[Install]
WantedBy=user-graphical-login.target

