# config/systemd/user/discord.service
[Service]
Type=exec
ExecStart=/usr/bin/discord

[Unit]
PartOf=graphical-session.target

[Install]
WantedBy=graphical-session.target
