# config/systemd/user/discord.service
[Service]
Type=exec
ExecStart=/usr/bin/discord-canary
[Install]
WantedBy=user-graphical-login.target
