# config/systemd/user/dunst.service
[Service]
ExecStart=/usr/bin/dunst
[Install]
WantedBy=user-graphical-login.target
