# config/systemd/user/dunst.service
[Service]
ExecStart=/usr/bin/sh -c 'echo "$XDG_CONFIG_HOME/dunst/dunstrc" | entr -nr dunst'
[Install]
WantedBy=user-graphical-login.target
