# config/systemd/user/dunst-watch.path
[Path]
PathChanged=%h/.config/dunst/dunstrc
[Install]
WantedBy=user-graphical-login.target
