# config/systemd/user/foot.service
[Service]
ExecStart=/usr/bin/foot --server
[Install]
WantedBy=user-graphical-login.target
