# config/systemd/user/foot.service
[Service]
ExecStart=/usr/bin/foot --server
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=user-graphical-login.target
