# config/systemd/user/generate-gtk-theme.service
[Service]
ExecStart=/usr/bin/sh -c 'generate depends gtk_theme | entr -np generate gtk_theme'
[Install]
WantedBy=user-graphical-login.target
