# config/systemd/user/generate-dunstrc.service
[Service]
ExecStart=/usr/bin/sh -c 'generate depends dunstrc | entr -np generate dunstrc'
[Install]
WantedBy=user-graphical-login.target
