# config/systemd/user/generate-xresources.service
[Service]
ExecStart=%h/.local/scriptp/generate xresources
[Install]
WantedBy=user-graphical-login.target
