# config/systemd/user/generate-thumbnails.service
[Service]
ExecStart=%h/.local/script/generate-thumbnails
