# config/systemd/user/generate-dunstrc.path
[Path]
PathChanged=%h/.local/script/generate
PathChanged=%h/.local/script/base16.sh
PathChanged=%h/.local/script/colors.sh
PathChanged=%h/.config/template/dunstrc
[Install]
WantedBy=user-graphical-login.target
