# config/systemd/user/generate-gtk-theme.service
[Service]
ExecStart=%h/.local/script/generate gtk_theme
[Install]
WantedBy=user-graphical-login.target
