# config/systemd/user/generate-dunstrc.service
[Service]
ExecStart=%h/.local/scriptp/generate dunstrc
[Install]
WantedBy=user-graphical-login.target
