# config/systemd/user/hide-cursor.service
[Service]
Type=exec
ExecStart=/usr/bin/unclutter --timeout 2
[Install]
WantedBy=user-graphical-login.target

