# config/systemd/user/hyprland-session.target
[Unit]
After=default.target
Requires=default.target
BindsTo=graphical-session.target
