# config/systemd/user/hyprland.service
[Unit]
BindsTo=graphical-session.target

[Service]
Type=notify
ExecStart=/usr/bin/Hyprland

[Install]
WantedBy=default.target
