# config/systemd/user/mail.service
[Service]
Type=exec
ExecStart=/usr/bin/thunderbird
[Install]
WantedBy=user-graphical-login.target

