# config/systemd/user/ollama.service
[Service]
Type=exec
ExecStart=/usr/bin/ollama serve
Environment="HSA_OVERRIDE_GFX_VERSION=10.3.0"

[Install]
WantedBy=user-graphical-login.target
