# config/systemd/user/optimize.service
[Service]
ExecStart=%h/.local/script/optimize %h/Pictures %h/Camera


