# config/systemd/user/picom.service
[Unit]
Description=Compositor

[Service]
ExecStart=/usr/bin/picom
Restart=always
RestartSec=500ms

[Install]
WantedBy=user-graphical-login.target

