# config/systemd/user/poweralert.service
[Service]
Type=exec
ExecStart=/usr/bin/poweralertd
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=user-graphical-login.target
