# config/systemd/user/rss-sync.service
[Service]
ExecStart=%h/.local/script/rss-sync
