# config/systemd/user/steam.service
[Service]
ExecStart=/usr/bin/steam
[Install]
WantedBy=user-graphical-login.target

