# config/systemd/user/steam.service
[Service]
Type=exec
ExecStart=%h/.local/script/steam
[Install]
WantedBy=user-graphical-login.target

