# config/systemd/user/steam.service
[Service]
Type=exec
ExecStart=%h/.local/script/steamwrap

[Install]
WantedBy=user-graphical-login.target

