# config/systemd/user/steam.service
[Service]
Type=exec
ExecStart=/usr/bin/steam -bigpicture

[Install]
WantedBy=user-graphical-login.target

