# config/systemd/user/sxhkd.service
[Service]
ExecStart=/usr/bin/sxhkd -r /tmp/sxhkd.log
[Install]
WantedBy=user-graphical-login.target
