# config/systemd/user/sxhkd-watch.path
[Path]
PathChanged=%h/.config/sxhkd/sxhkdrc
[Install]
WantedBy=user-graphical-login.target
