# config/systemd/user/sxhkd-watch.service
[Service]
Type=oneshot
ExecStart=/usr/bin/killall -s USR1 sxhkd

[Install]
WantedBy=user-graphical-login.target
