# config/systemd/user/sync-anime.service
[Unit]
Conflicts/usr/bin/shutdown.target
[Service]
ExecStart=/usr/local/bin/anilist fetch
