# config/systemd/user/sync-aniz.service
[Service]
ExecStart=%h/.local/bin/aniz database download


