# config/systemd/user/syncthing.service
[Service]
ExecStart=/usr/bin/syncthing -no-browser
[Install]
WantedBy=user-graphical-login.target
