# config/systemd/user/tmux-kill-detached-sessions.service
[Unit]
Conflicts=/usr/bin/shutdown.target
[Service]
ExecStart=%h/.local/script/tmux-kill-detatched-sessions
