# config/systemd/user/udiskie.service
[Service]
Type=exec
ExecStart=/usr/bin/udiskie
[Install]
WantedBy=user-graphical-login.target

