# config/systemd/user/udiskie.service
[Service]
Type=exec
ExecStart=/usr/bin/udiskie
Restart=on-failure
RestartSec=5s


[Install]
WantedBy=user-graphical-login.target

