# config/systemd/user/wallpaper.service
[Service]
ExecStart=%h/.local/script/wall.sh '%h/MEGA/images/wallpapers'
[Install]
WantedBy=user-graphical-login.target
