# config/systemd/user/wallpaper.service
[Service]
ExecStart=%h/.local/script/wall '%h/Sync/images/wallpapers'
[Install]
WantedBy=user-graphical-login.target
