# config/systemd/user/hyprpaper.service.d/generate.conf
[Service]
# Generate a per-monitor hyprpaper config (one wallpaper{} block per monitor)
# before hyprpaper starts, so each monitor gets a distinct random wallpaper.
# The main hyprpaper.conf `source`s this generated file.
ExecStartPre=%h/.local/script/gen-hyprpaper-conf
